# sys-apps/earlyoom/files/earlyoom
#!/sbin/openrc-run
name=earlyoom
command=/usr/bin/earlyoom
command_args=-p
command_background=1
pidfile=/run/earlyoom.pid
